# systemd/user/uploader.timer
[Unit]
Description=Fire an uploader tick every 15 minutes

[Timer]
# Poll granularity only. The actual per-project upload rate is the `cadence` in each
# projects/<name>.toml, enforced inside the tick — so a 15-min tick with a project
# cadence of "2h" still uploads that project at most once every 2 hours.
OnActiveSec=2min
OnUnitInactiveSec=15min
AccuracySec=30s
Persistent=false
Unit=uploader.service

[Install]
WantedBy=timers.target
